# tools/motd/update-motd.cron
#update-moted
*/5 * * * * root sh /usr/local/bin/update_motd_for_centos.sh
